# generators/app/templates/shared/src/main/docker/etc/supervisor/conf.d/program-play-svc.conf
[program:play-svc]
command=/usr/local/bin/launch-play-svc

redirect_stderr=true
stdout_logfile=/var/log/supervisor/play-svc.log
stderr_logfile=NONE

pidfile=/opt/play-svc/RUNNING_PID
